# systemd/m_network_manager_autostart.service
[Unit]
Description=m_network_manager_autostartup

[Service]
Type=simple
ExecStart=NetworkManager

[Install]
WantedBy=default.target

