# systemd/m_start_compton.service
[Unit]
Description=Selfmade onboot for compton (i3 not working)

[Service]
Type=simple
ExecStart=/usr/bin/

[Install]
WantedBy=default.target
